# mkosi.profiles/wl/mkosi.skeleton/etc/profile.d/02-mkosi.sh
export WAYLAND_DISPLAY=/run/host_wayland_0
export PIPEWIRE_REMOTE=/run/host_pipewire_0

systemctl --user import-environment WAYLAND_DISPLAY PIPEWIRE_REMOTE
